# pear-supervisor.conf
[program:pear]
command=bash /home/ubuntu/pear-admin-flask/start.sh
user=root
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true